# deploy/systemd/framenest.service
[Unit]
Description=FrameNest local-first media library server
Documentation=https://github.com/cisarik/framenest
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=framenest
Group=framenest
WorkingDirectory=/opt/framenest/current
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/framenest/framenest.env
ExecStartPre=/opt/framenest/current/.venv/bin/framenest-production check-database-ready
ExecStart=/opt/framenest/current/.venv/bin/framenest-production serve
KillSignal=SIGTERM
TimeoutStopSec=30s
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
StateDirectory=framenest
CacheDirectory=framenest
RuntimeDirectory=framenest
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
